# adaptive_2D_scaling/vertical_scaling/client/client-lib/ld.so.preload
/libcontroller.so
/libcuda.so
/libcuda.so.1
/libnvidia-ml.so
/libnvidia-ml.so.1
